# cooling-control.service
[Unit]
Description=Cooling Control — temperature monitor and fan/pump controller
After=multi-user.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/cooling-control/cooling-control
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/cooling-control

# Required capabilities for hardware access
CapabilityBoundingSet=CAP_SYS_RAWIO
AmbientCapabilities=CAP_SYS_RAWIO

[Install]
WantedBy=multi-user.target
